# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid

[program:stt-api]
command=/opt/stt-venv/bin/uvicorn app.main:app --app-dir /opt/app --host 127.0.0.1 --port 8001 --workers 1
directory=/opt/app
priority=10
autostart=true
autorestart=unexpected
startsecs=3
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:vllm]
command=/bin/bash -c 'exec vllm serve "$LLM_MODEL_ID" --host 127.0.0.1 --port 8000 --served-model-name exaone-yaho --quantization bitsandbytes --dtype bfloat16 --max-model-len "$VLLM_MAX_MODEL_LEN" --max-num-seqs "$VLLM_MAX_NUM_SEQS" --gpu-memory-utilization "$VLLM_GPU_MEMORY_UTILIZATION" --enable-prefix-caching --api-key "$API_KEY"'
priority=10
autostart=true
autorestart=unexpected
startsecs=10
stopasgroup=true
killasgroup=true
stopsignal=TERM
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
priority=20
autostart=true
autorestart=unexpected
startsecs=3
stopsignal=QUIT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
